# deploy/supervisor/laravel-worker.conf
; ═══════════════════════════════════════════════════════════════
; Supervisor Configuration — LUXE E-commerce Production
; ═══════════════════════════════════════════════════════════════
;
; This file manages two processes:
;   1. Queue worker — processes async jobs (emails, SMS, notifications, exports)
;   2. Socket server — Node.js Socket.IO bridge for real-time WebSocket events
;
; Prerequisites:
;   - Set QUEUE_CONNECTION=database in Laravel .env
;   - Run: php artisan queue:table && php artisan migrate
;   - Node.js 18+ installed for the socket server
;
; Deploy:
;   sudo cp this-file /etc/supervisor/conf.d/luxe.conf
;   sudo supervisorctl reread && sudo supervisorctl update
;   sudo supervisorctl start luxe-worker:* luxe-socket:*
;
; Queue Driver:  database (no Redis required)
; Logs:          storage/logs/worker.log, ../socket-server/logs/socket.log
;
; IMPORTANT:
;   - Laravel Horizon is NOT available (requires Redis)
;   - Replace /var/www/html with your actual project path
; ═══════════════════════════════════════════════════════════════

; ── Queue Worker (database driver) ──
; Processes async jobs (SendEmailJob, SendSmsJob, SendNotificationJob, etc.)
; Uses database driver — NO Redis required.
; Restarts every hour (--max-time=3600) to prevent memory leaks.
[program:luxe-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log
stopwaitsecs=3600

; ── Socket.IO Server (Node.js) ──
; Bridges real-time events from Laravel to connected React clients.
; Processes: order updates, notifications, chat messages, ad campaigns.
; Requires Node.js 18+ and npm install in the socket-server directory.
[program:luxe-socket]
process_name=%(program_name)s
command=node /var/www/html/socket-server/server.js
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/socket-server.log
stopwaitsecs=10
environment=NODE_ENV="production"
